# router/include/router.service
[Unit]
Description=Router daemon

[Service]
ExecStart=/usr/bin/router
Restart=always

[Install]
WantedBy=multi-user.target
