# builder/include/docker-system-prune.timer
[Unit]
Description=Docker System Prune Timer

[Timer]
OnUnitActiveSec=2h
Unit=docker-system-prune.service

[Install]
WantedBy=timers.target
