# builder/include/docker-system-prune.service
[Unit]
Description=Docker System Prune Service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -f
