# nas/include/zfs-scrub@.timer
[Unit]
Description=Monthly ZFS Pool Scrub Timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
